# copy_to-homedir_config_autostart/lxterminal-autostart.desktop
[Desktop Entry]
Encoding=UTF-8
Name=Terminal autostart
Comment=Start supercollider
Exec=/usr/bin/lxterminal -e "sclang"
